# systemd/bluex_export.timer
[Unit]
Description=Run Blue-X Audit Export daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
